# deploy/macos/com.usefulgitinfo.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.usefulgitinfo.watchdog</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__ROOT__/scripts/macos/watchdog.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__ROOT__</string>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>60</integer>

  <key>StandardOutPath</key>
  <string>__ROOT__/.runtime/launchd/watchdog.out.log</string>
  <key>StandardErrorPath</key>
  <string>__ROOT__/.runtime/launchd/watchdog.err.log</string>
</dict>
</plist>
